# utils/opencloud.service
[Unit]
Description=OpenCloud Backend Service
After=network.target

[Service]
ExecStart=/home/ubuntu/OpenCloud/bin/opencloud

# Restart automatically on crash or reboot
Restart=always
RestartSec=5

# Run as a non-root user (recommended)
User=ubuntu
Group=ubuntu

# Environment variables for your app
Environment=PORT=3030
Environment=ENV=production
Environment=LOG_LEVEL=info
# %U expands to the numeric UID for User=, so Podman uses that user's runtime dir.
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=CONTAINER_HOST=unix:/run/user/%U/podman/podman.sock

# Optional: redirect logs to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
